# systemd/feature-skills-webapp.service
[Unit]
Description=feature-skills webapp
After=default.target
# Crash-loop cap. StartLimit* are [Unit]-only keys in modern systemd; in
# [Service] they are silently ignored and the manager defaults (10s/5) apply.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
Environment=FEATURE_SKILLS_WEBAPP_PORT=8800
ExecStart=%h/.local/bin/feature-skills-webapp
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
